# deploy/legacy/com.seongqkim.economy-agent.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.seongqkim.economy-agent</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>/Users/barq/developer/projects/economy-content-agent/run_agent.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/barq/developer/projects/economy-content-agent</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/economy-agent.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/economy-agent.err.log</string>
  </dict>
</plist>
